# deploy/berserk-tracker.service
[Unit]
Description=Berserk Manga Availability Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=berserk-tracker
Group=berserk-tracker
WorkingDirectory=/opt/berserk-tracker
ExecStart=/opt/berserk-tracker/venv/bin/python /opt/berserk-tracker/berserk_tracker.py
Restart=always
RestartSec=30

# Variaveis de ambiente
EnvironmentFile=/etc/berserk-tracker/config.env

# Seguranca
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/berserk-tracker
PrivateTmp=true

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=berserk-tracker

[Install]
WantedBy=multi-user.target
